# coolwater.conf
[program:coolwater]
command=/usr/local/bin/gunicorn wsgi --bind localhost:2000 --workers 3 --worker-class gevent
directory=/var/www/coolwater
autostart=true
autorestart=true
